# packaging/packetsonded.service
[Unit]
Description=packetsonde network observation agent
Documentation=https://github.com/billndotnet/packetsonde
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/packetsonded -c /etc/packetsonded/packetsonded.toml
Restart=on-failure
RestartSec=5s

# Capabilities: pcap requires NET_RAW + NET_ADMIN; mTLS doesn't need them
# itself but the agent's passive modules do.
# NOTE: enabling the [detect] track (fanotify file/exec monitoring) additionally
# requires CAP_SYS_ADMIN (fanotify_init) and CAP_DAC_READ_SEARCH (inode marks on
# unreadable sensitive files like /etc/shadow):
#   AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
#   CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Run as an unprivileged user with the right capabilities granted by
# systemd. Create with: useradd --system --no-create-home --shell /usr/sbin/nologin packetsonded
User=packetsonded
Group=packetsonded

# Filesystem isolation. The agent only needs read access to its config
# and key dirs, plus write to its state/log dirs.
ReadOnlyPaths=/etc/packetsonded
# Must match the [detect] sink/baseline dir (/var/lib/packetsonde) so it stays
# writable under ProtectSystem=strict.
StateDirectory=packetsonde
LogsDirectory=packetsonde
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
